# christmas.service
[Unit]
Description=Christmas lights
After=getty.target

[Service]
ExecStart=/home/tom/christmas_start.sh

[Install]
WantedBy=default.target

